# crontab
#LOTTO
05 21 * * 6-7 /home/jihwan/lotto/checker.py &> /dev/null 2>&1
00 18 * * 5 /home/jihwan/lotto/reminder.py &> /dev/null 2>&1
