# crontab
#LOTTO
05 21 * * 6-7 ~/lotto/checker.py &> /dev/null 2>&1
00 18 * * 5 ~/lotto/reminder.py &> /dev/null 2>&1
